# ympd/templates/ympd.conf.upstart.j2
description     "ympd"

start on local-filesystems or runlevel [345]
stop on runlevel [!345]


exec /usr/local/bin/ympd -u {{ ympd_user }} -w {{ ympd_listen_port }}
