# lirc/files/upstart-ir-keytable.conf
# ir-keytable
#
# This task restores ir-keytable mappings

description	"Restore infrared remote keytables"

start on runlevel [2345]

task

exec     ir-keytable -a /etc/rc_maps.cfg
