# jupyterhub/templates/jupyterhub.service.j2
[Unit]
Description=JupyterHub
After=network.target

[Service]
WorkingDirectory=/opt/jupyterhub/
ExecStart=/opt/jupyter/venv/bin/jupyterhub

[Install]
WantedBy=multi-user.target
